# kiwix-rag.service
# ── Edit these to match your setup ──────────────────────────────────────────
# User/Group:   your Pi username
# ExecStart:    path to your venv Python and the script on your SSD
# --db:         ChromaDB directory
# --embed-model: directory containing the all-MiniLM-L6-v2 snapshot
# --model:      Ollama model name (run `ollama list` to see installed models)
# RequiresMountsFor: your SSD mount point (or remove if running from SD card)
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Kiwix RAG Web Interface
After=network.target ollama.service
Requires=ollama.service
RequiresMountsFor=/mnt/ssd

[Service]
Type=simple
User=pi
Group=pi
Environment="HF_HUB_OFFLINE=1"
ExecStart=/home/pi/kiwix-rag/bin/python /mnt/ssd/kiwix-rag-project/web.py \
    --host 0.0.0.0 \
    --db /mnt/ssd/vector_db \
    --embed-model /mnt/nvme/all-MiniLM-L6-v2 \
    --ollama-url http://cori-desktop.local:11434 \
    --model llama3.1:8b-rag \
    --max-cache-bytes 11000000000 \
    --max-per-group 5
Restart=on-failure
RestartSec=15
# Collection HNSW indexes are 2-7 GB each; cache is byte-bounded so total
# resident index memory stays within budget. On a 16 GB Pi, cap resident memory
# so an overflow kills only this service (clean restart) instead of triggering
# a global OOM that takes down networking. Pair with --max-cache-bytes above.
MemoryMax=13G
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kiwix-rag

[Install]
WantedBy=multi-user.target
